# ansible/roles/geoipupdate/templates/geoipupdate.cron.j2
#!/bin/sh
# {{ ansible_managed }}
#

set -e

/usr/bin/geoipupdate
