# hacking/platforms/ubuntu/cif-router.cron
#!/bin/sh

set -e

/etc/init.d/cif-router restart