# hacking/platforms/ubuntu/cif-router.cron
#!/bin/bash

/etc/init.d/cif-router restart