# hacking/platforms/ubuntu/geoipupdate.cron
#!/bin/sh

set -e

/usr/bin/geoipupdate
