# ansible/roles/cifv2/templates/cif-router.cron.j2
#!/bin/sh
# {{ ansible_managed }}
#

set -e

/etc/init.d/cif-router restart > /dev/null 2>&1
